# ops/systemd/gateway-blue.service
# =============================================================================
# ops/systemd/gateway-blue.service – systemd unit for the blue gateway instance
#
# INSTALLATION
#   sudo cp ops/systemd/gateway-blue.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable gateway-blue.service
#
# The unit reads secrets from the environment file at:
#   /etc/home-automation/gateway-blue.env
# Copy ops/systemd/gateway-blue.env.example to that path and fill in values.
# =============================================================================

[Unit]
Description=Home Automation Gateway – Blue instance
Documentation=file:///opt/home-automation-scripts/docs/blue_green.md
After=network.target
# Allow the green instance to start independently.
Wants=network.target

[Service]
Type=simple
User=home-automation
Group=home-automation

# Load instance-specific secrets and config.
# This file must exist before starting the service; see *.env.example.
EnvironmentFile=/etc/home-automation/gateway-blue.env

WorkingDirectory=/opt/home-automation-scripts

# Use gunicorn for a production-grade WSGI server.
# Workers=1 is intentional for a single-process, low-traffic home server.
ExecStart=/opt/home-automation-scripts/.venv/bin/gunicorn \
    --workers 1 \
    --bind 127.0.0.1:${GATEWAY_PORT} \
    --access-logfile - \
    --error-logfile - \
    --capture-output \
    gateway.app:app

# Restart on failure but not on clean exit.
Restart=on-failure
RestartSec=5s

# Give the process 30 seconds to shut down cleanly before SIGKILL.
TimeoutStopSec=30

# Security hardening – appropriate for a trusted home-server environment.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/home-automation-scripts/logs

[Install]
WantedBy=multi-user.target
