# ops/systemd/gateway-green.service
# =============================================================================
# ops/systemd/gateway-green.service – systemd unit for the green gateway instance
#
# INSTALLATION
#   sudo cp ops/systemd/gateway-green.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable gateway-green.service
#
# The unit reads secrets from the environment file at:
#   /etc/home-automation/gateway-green.env
# Copy ops/systemd/gateway-green.env.example to that path and fill in values.
# =============================================================================

[Unit]
Description=Home Automation Gateway – Green instance
Documentation=file:///opt/home-automation-scripts/docs/blue_green.md
After=docker.service network.target
Requires=docker.service
Wants=network.target

[Service]
Type=simple
User=home-automation
Group=home-automation
SupplementaryGroups=docker

# Load instance-specific secrets and config.
EnvironmentFile=/etc/home-automation/gateway-green.env

ExecStartPre=-/usr/bin/docker rm -f gateway-green
ExecStart=/bin/sh -lc '/usr/bin/docker run --rm --name gateway-green --env-file /etc/home-automation/gateway-green.env --publish 127.0.0.1:${GATEWAY_PORT}:${GATEWAY_PORT} "${GATEWAY_IMAGE}"'
ExecStop=-/usr/bin/docker stop gateway-green

Restart=on-failure
RestartSec=5s

TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/run/docker.sock

[Install]
WantedBy=multi-user.target
