# ops/systemd/gateway-blue.service
# =============================================================================
# ops/systemd/gateway-blue.service – systemd unit for the blue gateway instance
#
# INSTALLATION
#   sudo cp ops/systemd/gateway-blue.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable gateway-blue.service
#
# The unit reads secrets from the environment file at:
#   /etc/home-automation/gateway-blue.env
# Copy ops/systemd/gateway-blue.env.example to that path and fill in values.
# =============================================================================

[Unit]
Description=Home Automation Gateway – Blue instance
Documentation=file:///opt/home-automation-scripts/docs/blue_green.md
After=docker.service network.target
Requires=docker.service
# Allow the green instance to start independently.
Wants=network.target

[Service]
Type=simple
User=home-automation
Group=home-automation
SupplementaryGroups=docker

# Load instance-specific secrets and config.
# This file must exist before starting the service; see *.env.example.
# It also declares the Docker image tag for this slot.
EnvironmentFile=/etc/home-automation/gateway-blue.env

ExecStartPre=-/usr/bin/docker rm -f gateway-blue
ExecStart=/bin/sh -lc '/usr/bin/docker run --rm --name gateway-blue --env-file /etc/home-automation/gateway-blue.env --publish 127.0.0.1:${GATEWAY_PORT}:${GATEWAY_PORT} "${GATEWAY_IMAGE}"'
ExecStop=-/usr/bin/docker stop gateway-blue

# Restart on failure but not on clean exit.
Restart=on-failure
RestartSec=5s

# Give the process 30 seconds to shut down cleanly before SIGKILL.
TimeoutStopSec=30

# Security hardening – appropriate for a trusted home-server environment.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/run/docker.sock

[Install]
WantedBy=multi-user.target
