# ops/systemd/timers/home-automation-package-update-notifier.service
[Unit]
Description=Home Automation Package Update Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=home-automation
Group=home-automation
EnvironmentFile=-/etc/home-automation/automation.env
WorkingDirectory=__REPO_ROOT__
ExecStart=__PYTHON__ __REPO_ROOT__/tools/automation.py run package-update-notifier

