# ops/systemd/timers/home-automation-health-check.timer
[Unit]
Description=Run home automation health check every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=home-automation-health-check.service

[Install]
WantedBy=timers.target

